# vm-config/init/promo-pay.conf
# Upstart script for the Promo Pay server

description "Promo Pay Server"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown
respawn

setuid promo-pay
setgid promo-pay

script
  export NODE_ENV="production"
  exec promo-pay >> /var/opt/promo-pay/log/all.log 2>&1
end script
